# bootstrap/supervisord.conf
[supervisord]
nodaemon=true

[program:redis_cluster]
command=/home/redis/redis/bootstrap/cluster.sh
user=redis
startsecs=1
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
